# setup/openebsdriver.socket
[Unit]
Description=openebsdriver

[Socket]
ListenStream=/run/docker/plugins/openebsdriver.sock

[Install]
WantedBy=sockets.target

